# systemd/openbuild.service
[Unit]
Description=Run openbuild build

[Service]
Type=oneshot
ExecStart=/usr/bin/openbuild-run
User=openbuild
